# systemd/clash-socks-webui.service
[Unit]
Description=Clash Socks Server UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/clash-socks-server-ui
EnvironmentFile=/opt/clash-socks-server-ui/.env
ExecStart=/opt/clash-socks-server-ui/runtime/start-webui.sh
Restart=always
RestartSec=3
UMask=027
KillMode=control-group
TimeoutStopSec=10
User=clashui
Group=clashui

[Install]
WantedBy=multi-user.target
